# supervisord.conf
[supervisord]
nodaemon = true
loglevel = debug
directory = /tmp

[program:configure]
command = /configure-script
priority = 1
startsecs = 0
exitcodes = 0

[program:php-fpm]
command = /usr/sbin/php-fpm7.0 --nodaemonize -c /etc/php/7.0/fpm
priority = 2
stdout_logfile = /dev/fd/1
redirect_stderr = true

[program:movim]
priority = 3
command = /usr/bin/php /var/www/html/daemon.php %(ENV_MOVIM_URL)s %(ENV_MOVIM_PORT)s
stdout_logfilei = /dev/fd/1
redirect_stderr = true
user = www-data

[program:nginx]
priority = 4
command = /usr/sbin/nginx -g "daemon off;"
stdout_logfile = /dev/fd/1
redirect_stderr = true

